# crontab
0 5 * * * /mnt/projects/eyerisserver/releases/current/eyerisserver/jobs.sh
